# systemd/cooler-lcd.service
[Unit]
Description=Cooler LCD sensor feed (Vevor/HT HID 5131:2007)
After=graphical-session.target

[Service]
ExecStart=/usr/bin/python3 %h/vevor_lcd_linux.py
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
